# misc/lala-musicbar.service.in
[Unit]
Description=A music bar with iced
PartOf=graphical-session.target
After=graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
ExecStart=@bin_install_dir@/lala_musicbar
Restart=on-failure

[Install]
WantedBy=graphical-session.target
